# packaging/systemd/nut-slave.service
[Unit]
Description=nut slave service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nut-server
Group=nut-server
WorkingDirectory=/usr/local/lib/nut-server
ExecStart=/usr/local/bin/nut-slave -config /etc/nut-server/slave.yaml
Restart=always
RestartSec=5

ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
ReadWritePaths=/var/lib/nut-server
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
UMask=0027

[Install]
WantedBy=multi-user.target
