# packaging/systemd/nut-master.service
[Unit]
Description=nut master service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/lib/nut-server
ExecStart=/usr/local/bin/nut-master -config /etc/nut-server/master.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
